# chair/config/templates/systemd/vmraid-chair-web.target
[Unit]
After=network.target
Wants={{ chair_name }}-vmraid-web.service {{ chair_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
