# chair/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:{{ chair_name }}-vmraid-web]
command={{ chair_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} -t {{ http_timeout }} vmraid.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/web.log
stderr_logfile={{ chair_dir }}/logs/web.error.log
user={{ user }}
directory={{ sites_dir }}

{% if use_rq %}
[program:{{ chair_name }}-vmraid-schedule]
command={{ chair_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/schedule.log
stderr_logfile={{ chair_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ chair_dir }}

[program:{{ chair_name }}-vmraid-default-worker]
command={{ chair_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ chair_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ chair_name }}-vmraid-short-worker]
command={{ chair_cmd }} worker --queue short
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ chair_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ chair_name }}-vmraid-long-worker]
command={{ chair_cmd }} worker --queue long
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ chair_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

{% for worker_name, worker_details in workers.items() %}
[program:{{ chair_name }}-vmraid-{{ worker_name }}-worker]
command={{ chair_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ chair_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endfor %}

{% else %}
[program:{{ chair_name }}-vmraid-workerbeat]
command={{ chair_dir }}/env/bin/python -m vmraid.celery_app beat -s beat.schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/workerbeat.log
stderr_logfile={{ chair_dir }}/logs/workerbeat.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ chair_name }}-vmraid-worker]
command={{ chair_dir }}/env/bin/python -m vmraid.celery_app worker -n jobs@%%h -Ofair --soft-time-limit 360 --time-limit 390 --loglevel INFO
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=400
directory={{ sites_dir }}
killasgroup=true

[program:{{ chair_name }}-vmraid-longjob-worker]
command={{ chair_dir }}/env/bin/python -m vmraid.celery_app worker -n longjobs@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

[program:{{ chair_name }}-vmraid-async-worker]
command={{ chair_dir }}/env/bin/python -m vmraid.celery_app worker -n async@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/worker.log
stderr_logfile={{ chair_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

{% endif %}

{% if not skip_redis %}
[program:{{ chair_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/redis-cache.log
stderr_logfile={{ chair_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ chair_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/redis-queue.log
stderr_logfile={{ chair_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if not skip_redis %}
[program:{{ chair_name }}-redis-socketio]
command={{ redis_server }} {{ redis_socketio_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/redis-socketio.log
stderr_logfile={{ chair_dir }}/logs/redis-socketio.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if node %}
[program:{{ chair_name }}-node-socketio]
command={{ node }} {{ chair_dir }}/apps/vmraid/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ chair_dir }}/logs/node-socketio.log
stderr_logfile={{ chair_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ chair_dir }}
{% endif %}

[group:{{ chair_name }}-web]
programs={{ chair_name }}-vmraid-web {%- if node -%} ,{{ chair_name }}-node-socketio {%- endif%}

{% if use_rq %}

[group:{{ chair_name }}-workers]
programs={{ chair_name }}-vmraid-schedule,{{ chair_name }}-vmraid-default-worker,{{ chair_name }}-vmraid-short-worker,{{ chair_name }}-vmraid-long-worker

{% else %}

[group:{{ chair_name }}-workers]
programs={{ chair_name }}-vmraid-workerbeat,{{ chair_name }}-vmraid-worker,{{ chair_name }}-vmraid-longjob-worker,{{ chair_name }}-vmraid-async-worker

{% endif %}

{% if not skip_redis %}
[group:{{ chair_name }}-redis]
programs={{ chair_name }}-redis-cache,{{ chair_name }}-redis-queue,{{ chair_name }}-redis-socketio
{% endif %}
